# configs/AppTemplate.service
[Unit]
Description=ExerciseDiary
Documentation=https://github.com/aceberg/ExerciseDiary
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/ExerciseDiary -c /etc/ExerciseDiary/config.yaml
Restart=on-failure

[Install]
WantedBy=multi-user.target